# fuel/app/tasks/cron.d/fleamarket_event_close
5 0 * * * apache php /var/www/html/oil refine fleamarket_event_close

