# vagrant/site-cookbooks/rakuichi-rakuza/templates/default/rakuichi-rakuza.crontab.erb
FUEL_ENV=<%= @mode %>

# フリマキャンセル待ちユーザ一への通知メール
#0 8 * * * apache (cd /deploy/rakuichi-rakuza ; /usr/bin/php oil refine send_mail_to_user_waiting_reservation)

# 予約受付開始日を過ぎたフリマの開催状況(event_status)を予約受付中に更新する
1 0 * * * apache php /deploy/rakuichi-rakuza/oil refine fleamarket_event_reservation_receipt

# 開催日を過ぎたフリマの開催状況(event_status)を開催終了に更新する
5 0 * * * apache php /deploy/rakuichi-rakuza/oil refine fleamarket_event_close
